# ai-services/inference/idle-watchdog.timer
[Unit]
Description=Run the Gemma idle watchdog every 3 minutes

[Timer]
OnBootSec=6min
OnUnitActiveSec=3min
AccuracySec=30s

[Install]
WantedBy=timers.target
